# config/systemd/user/ibus.service
[Unit]
Description=IBUS Daemon

[Service]
ExecStart=/usr/bin/ibus-daemon

[Install]
WantedBy=graphical.target
