# config/systemd/user/muchsync@.service
[Unit]
Description=Sync email with mailserver

[Service]
ExecStart=/usr/bin/env muchsync %i
