# config/systemd/user/muchsync@.timer
[Unit]
Description=Periodically fetch mail

[Timer]
OnStartupSec=10s
OnUnitActiveSec=1m
Unit=muchsync@%i.service

[Install]
WantedBy=default.target
