# config/systemd/user/muchsync@.timer
[Unit]
Description=Periodically fetch mail

[Timer]
OnActiveSec=1m
OnUnitActiveSec=1m
Unit=muchsync@%i.service

[Install]
WantedBy=timers.target
