# config/systemd/user/psd-resync.timer.d/interval.conf
[Timer]
OnUnitActiveSec=
OnUnitActiveSec=10m
